# nix-darwin/result/Library/LaunchDaemons/org.nixos.activate-system.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.nixos.activate-system</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/sh</string>
		<string>-c</string>
		<string>/bin/wait4path /nix/store &amp;&amp; exec /nix/store/3mz9k8fg8yiscfbdi4ip7kbv4a5p1s7z-activate-system-start</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>